# utils/systemd/update-repo-notify.service
[Unit]
Description=Update repo notify service

[Service]
ExecStart=/var/lib/update/trigger.py